# buildroot-hearth/board/hearth-pi5/overlay/etc/systemd/system/hearth.service
[Unit]
Description=Hearth Smart Home Kiosk
After=network-online.target systemd-modules-load.service depmod-firstboot.service
Wants=network-online.target
OnFailure=hearth-rollback.service

[Service]
Type=simple
User=root
ExecStart=/opt/hearth/bundle/flutter-pi --release /opt/hearth/bundle
Environment=LD_LIBRARY_PATH=/opt/hearth/bundle
Environment=HEARTH_NO_MEDIAKIT=1
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3
Environment=XDG_RUNTIME_DIR=/run/user/0

[Install]
WantedBy=multi-user.target
