# buildroot-hearth/board/hearth-pi5/overlay/etc/systemd/system/wyoming-openwakeword.service
[Unit]
Description=Wyoming openWakeWord
After=network.target

[Service]
Type=simple
User=hearth
ExecStart=/opt/wyoming/wyoming-openwakeword/.venv/bin/python3 -m wyoming_openwakeword \
    --uri tcp://127.0.0.1:10400 \
    --preload-model ok_nabu \
    --threshold 0.3
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
